# deploy/mikrotik-secrets-sync.timer
[Unit]
Description=Hourly IPDR MikroTik secrets sync

[Timer]
# First run 3 min after boot, then hourly. Randomized delay spreads load if
# multiple NAS are ever added.
OnBootSec=3min
OnUnitActiveSec=1h
RandomizedDelaySec=2min
Persistent=true

[Install]
WantedBy=timers.target
